# .autonoma/cron/com.flitterbot.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.flitterbot.scheduler</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__FLITTERBOT_DIR__/cron/flitterbot-checkin.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>StandardOutPath</key>
  <string>__FLITTERBOT_DIR__/logs/cron.log</string>
  <key>StandardErrorPath</key>
  <string>__FLITTERBOT_DIR__/logs/cron.log</string>
  <key>WorkingDirectory</key>
  <string>__HOME__</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
